# OK8MP-linux-fs/rootfs/var/lib/dpkg/info/alsa-utils-alsactl.list
/.
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/alsa-restore.service
/lib/systemd/system/alsa-state.service
/lib/systemd/system/sound.target.wants
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/90-alsa-restore.rules
/usr
/usr/sbin
/usr/sbin/alsactl
/usr/share
/usr/share/alsa
/usr/share/alsa/init
/usr/share/alsa/init/00main
/usr/share/alsa/init/ca0106
/usr/share/alsa/init/default
/usr/share/alsa/init/hda
/usr/share/alsa/init/help
/usr/share/alsa/init/info
/usr/share/alsa/init/test
/var
/var/lib
/var/lib/alsa
/lib/systemd/system/sound.target.wants/alsa-restore.service
/lib/systemd/system/sound.target.wants/alsa-state.service
