# OK8MP-linux-fs/rootfs/var/lib/dpkg/info/at.list
/.
/etc
/etc/at.deny
/etc/pam.d
/etc/pam.d/atd
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/atd.service
/lib/systemd/system-preset
/lib/systemd/system-preset/98-at.preset
/usr
/usr/bin
/usr/bin/at
/usr/bin/batch
/usr/sbin
/usr/sbin/atd
/usr/sbin/atrun
/var
/var/spool
/var/spool/at
/var/spool/at/jobs
/var/spool/at/jobs/.SEQ
/var/spool/at/spool
/usr/bin/atq
/usr/bin/atrm
