# OK8MP-linux-fs/rootfs/var/lib/dpkg/info/at.prerm
#!/bin/sh
[ "$1" != "upgrade" ] || exit 0
set -e
if type systemctl >/dev/null 2>/dev/null; then
	if [ -z "$D" ]; then
		systemctl stop atd.service

		systemctl disable atd.service
	fi
fi
