# OK8MP-linux-fs/rootfs/var/lib/dpkg/info/avahi-daemon.list
/.
/etc
/etc/avahi
/etc/avahi/avahi-daemon.conf
/etc/avahi/hosts
/etc/avahi/services
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/avahi-dbus.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/avahi-daemon.service
/lib/systemd/system/avahi-daemon.socket
/lib/systemd/system-preset
/lib/systemd/system-preset/98-avahi-daemon.preset
/usr
/usr/sbin
/usr/sbin/avahi-daemon
/usr/share
/usr/share/avahi
/usr/share/avahi/avahi-service.dtd
/usr/share/dbus-1
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/org.freedesktop.Avahi.service
