# OK8MP-linux-fs/rootfs/var/lib/dpkg/info/connman.prerm
#!/bin/sh
[ "$1" != "upgrade" ] || exit 0
set -e
	update-alternatives --remove  resolv-conf /etc/resolv-conf.connman
if type systemctl >/dev/null 2>/dev/null; then
	if [ -z "$D" ]; then
		systemctl stop connman.service

		systemctl disable connman.service
	fi
fi
