# OK8MP-linux-fs/rootfs/var/lib/dpkg/info/connman.list
/.
/etc
/etc/tmpfiles.d
/etc/tmpfiles.d/connman_resolvconf.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/connman.service
/lib/systemd/system-preset
/lib/systemd/system-preset/98-connman.preset
/usr
/usr/lib
/usr/lib/connman
/usr/lib/connman/plugins
/usr/sbin
/usr/sbin/connmand
/usr/share
/usr/share/dbus-1
/usr/share/dbus-1/system.d
/usr/share/dbus-1/system.d/connman.conf
/etc/resolv-conf.connman
