# OK8MP-linux-fs/rootfs/var/lib/dpkg/info/dbus-1.list
/.
/etc
/etc/dbus-1
/etc/dbus-1/session.conf
/etc/dbus-1/system.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/dbus.service
/lib/systemd/system/dbus.socket
/lib/systemd/system/dbus.target.wants
/lib/systemd/system/multi-user.target.wants
/lib/systemd/system/sockets.target.wants
/usr
/usr/bin
/usr/bin/dbus-cleanup-sockets
/usr/bin/dbus-daemon
/usr/bin/dbus-launch
/usr/bin/dbus-monitor
/usr/bin/dbus-run-session
/usr/bin/dbus-send
/usr/bin/dbus-update-activation-environment
/usr/bin/dbus-uuidgen
/usr/lib
/usr/lib/systemd
/usr/lib/systemd/user
/usr/lib/systemd/user/dbus.service
/usr/lib/systemd/user/dbus.socket
/usr/lib/systemd/user/sockets.target.wants
/usr/lib/sysusers.d
/usr/lib/sysusers.d/dbus.conf
/usr/lib/tmpfiles.d
/usr/lib/tmpfiles.d/dbus.conf
/usr/libexec
/usr/libexec/dbus-daemon-launch-helper
/usr/share
/usr/share/dbus-1
/usr/share/dbus-1/services
/usr/share/dbus-1/session.conf
/usr/share/dbus-1/session.d
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system.conf
/usr/share/dbus-1/system.d
/usr/share/xml
/usr/share/xml/dbus-1
/usr/share/xml/dbus-1/busconfig.dtd
/usr/share/xml/dbus-1/introspect.dtd
/var
/var/lib
/var/lib/dbus
/lib/systemd/system/dbus.target.wants/dbus.socket
/lib/systemd/system/multi-user.target.wants/dbus.service
/lib/systemd/system/sockets.target.wants/dbus.socket
/usr/lib/systemd/user/sockets.target.wants/dbus.socket
