# OK8MP-linux-fs/rootfs/var/lib/dpkg/info/dropbear.list
/.
/etc
/etc/default
/etc/default/dropbear
/etc/dropbear
/etc/pam.d
/etc/pam.d/dropbear
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/dropbear.socket
/lib/systemd/system/dropbear@.service
/lib/systemd/system/dropbearkey.service
/lib/systemd/system-preset
/lib/systemd/system-preset/98-dropbear.preset
/usr
/usr/bin
/usr/sbin
/usr/sbin/dropbearmulti
/var
/usr/bin/dbclient
/usr/sbin/dropbear
/usr/sbin/dropbearconvert
/usr/sbin/dropbearkey
