# OK8MP-linux-fs/rootfs/var/lib/dpkg/info/fsl-rc-local.list
/.
/etc
/etc/init.d
/etc/init.d/rc.local
/etc/rc.local
