# OK8MP-linux-fs/rootfs/var/lib/dpkg/info/hostapd.list
/.
/etc
/etc/hostapd.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/hostapd.service
/lib/systemd/system-preset
/lib/systemd/system-preset/98-hostapd.preset
/usr
/usr/sbin
/usr/sbin/hostapd
/usr/sbin/hostapd_cli
