# OK8MP-linux-fs/rootfs/var/lib/dpkg/info/iptables.list
/.
/etc
/etc/iptables
/etc/iptables/iptables.rules
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/iptables.service
/lib/systemd/system-preset
/lib/systemd/system-preset/98-iptables.preset
/usr
/usr/bin
/usr/lib
/usr/lib/libip4tc.so.2.0.0
/usr/lib/libip6tc.so.2.0.0
/usr/lib/libiptc.so.0.0.0
/usr/lib/libxtables.so.12.2.0
/usr/sbin
/usr/sbin/xtables-legacy-multi
/usr/bin/iptables-xml
/usr/lib/libip4tc.so.2
/usr/lib/libip6tc.so.2
/usr/lib/libiptc.so.0
/usr/lib/libxtables.so.12
/usr/sbin/ip6tables
/usr/sbin/ip6tables-legacy
/usr/sbin/ip6tables-legacy-restore
/usr/sbin/ip6tables-legacy-save
/usr/sbin/ip6tables-restore
/usr/sbin/ip6tables-save
/usr/sbin/iptables
/usr/sbin/iptables-legacy
/usr/sbin/iptables-legacy-restore
/usr/sbin/iptables-legacy-save
/usr/sbin/iptables-restore
/usr/sbin/iptables-save
