# OK8MP-linux-fs/rootfs/var/lib/dpkg/info/iputils-ninfod.list
/.
/etc
/etc/init.d
/etc/init.d/ninfod.sh
/sbin
/sbin/ninfod
