# OK8MP-linux-fs/rootfs/var/lib/dpkg/info/isp-imx.list
/.
/etc
/etc/default
/etc/default/imx8-isp
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/imx8-isp.service
/lib/systemd/system-preset
/lib/systemd/system-preset/98-isp-imx.preset
/opt
/opt/imx8-isp
/opt/imx8-isp/bin
/opt/imx8-isp/bin/case
/opt/imx8-isp/bin/case/config
/opt/imx8-isp/bin/isp_media_server
/opt/imx8-isp/bin/run.sh
/opt/imx8-isp/bin/start_isp.sh
/usr
/usr/lib
/usr/lib/liba2dnr.so
/usr/lib/liba3dnr.so
/usr/lib/libadpcc.so
/usr/lib/libadpf.so
/usr/lib/libaec.so
/usr/lib/libaee.so
/usr/lib/libaf.so
/usr/lib/libaflt.so
/usr/lib/libahdr.so
/usr/lib/libappshell_ebase.so
/usr/lib/libappshell_hal.so
/usr/lib/libappshell_ibd.so
/usr/lib/libappshell_oslayer.so
/usr/lib/libavs.so
/usr/lib/libawb.so
/usr/lib/libawdr3.so
/usr/lib/libbase64.so
/usr/lib/libbufferpool.so
/usr/lib/libbufsync_ctrl.so
/usr/lib/libcam_calibdb.so
/usr/lib/libcam_device.so
/usr/lib/libcam_engine.so
/usr/lib/libcameric_drv.so
/usr/lib/libcameric_reg_drv.so
/usr/lib/libcim_ctrl.so
/usr/lib/libcommon.so
/usr/lib/libdewarp_hal.so
/usr/lib/libebase.so
/usr/lib/libfpga.so
/usr/lib/libhal.so
/usr/lib/libi2c_drv.so
/usr/lib/libibd.so
/usr/lib/libisi.so
/usr/lib/libjsoncpp.so.1.9.0
/usr/lib/libmedia_server.so
/usr/lib/libmim_ctrl.so
/usr/lib/libmipi_drv.so
/usr/lib/libmom_ctrl.so
/usr/lib/libos08a20.so.1.0.0
/usr/lib/liboslayer.so
/usr/lib/libov2775.so.1.0.0
/usr/lib/libsom_ctrl.so
/usr/lib/libtinyxml2.so
/usr/lib/libversion.so
/usr/lib/libvom_ctrl.so
/usr/lib/libvvdisplay_shared.so
/usr/lib/libjsoncpp.so
/usr/lib/libjsoncpp.so.21
/usr/lib/libos08a20.so
/usr/lib/libos08a20.so.1
/usr/lib/libov2775.so
/usr/lib/libov2775.so.1
