# OK8MP-linux-fs/rootfs/var/lib/dpkg/info/isp-imx.prerm
#!/bin/sh
[ "$1" != "upgrade" ] || exit 0
set -e
if type systemctl >/dev/null 2>/dev/null; then
	if [ -z "$D" ]; then
		systemctl stop imx8-isp.service

		systemctl disable imx8-isp.service
	fi
fi
