# OK8MP-linux-fs/rootfs/var/lib/dpkg/info/lighttpd.list
/.
/etc
/etc/lighttpd
/etc/lighttpd/lighttpd.conf
/etc/lighttpd.d
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/lighttpd.service
/lib/systemd/system-preset
/lib/systemd/system-preset/98-lighttpd.preset
/usr
/usr/sbin
/usr/sbin/lighttpd
/usr/sbin/lighttpd-angel
/www
/www/pages
/www/pages/dav
/www/pages/index.html
/www/logs
/www/var
