# OK8MP-linux-fs/rootfs/var/lib/dpkg/info/logrotate.list
/.
/etc
/etc/logrotate.conf
/etc/logrotate.d
/etc/logrotate.d/btmp
/etc/logrotate.d/wtmp
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/logrotate.service
/lib/systemd/system/logrotate.timer
/lib/systemd/system-preset
/lib/systemd/system-preset/98-logrotate.preset
/usr
/usr/sbin
/usr/sbin/logrotate
/var
/var/lib
/var/lib/logrotate.status
