# OK8MP-linux-fs/rootfs/var/lib/dpkg/info/lttng-tools.list
/.
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/lttng-sessiond.service
/lib/systemd/system-preset
/lib/systemd/system-preset/98-lttng-tools.preset
/usr
/usr/bin
/usr/bin/lttng
/usr/bin/lttng-crash
/usr/bin/lttng-relayd
/usr/bin/lttng-sessiond
/usr/lib
/usr/lib/liblttng-ctl.so.0.0.0
/usr/lib/lttng
/usr/lib/lttng/libexec
/usr/lib/lttng/libexec/lttng-consumerd
/usr/share
/usr/share/xml
/usr/share/xml/lttng
/usr/share/xml/lttng/session.xsd
/usr/lib/liblttng-ctl.so.0
