# OK8MP-linux-fs/rootfs/var/lib/dpkg/info/lvm2.prerm
#!/bin/sh
[ "$1" != "upgrade" ] || exit 0
set -e
if type systemctl >/dev/null 2>/dev/null; then
	if [ -z "$D" ]; then
		systemctl stop lvm2-monitor.service dm-event.socket dm-event.service blk-availability.service lvm2-pvscan@.service

		systemctl disable lvm2-monitor.service dm-event.socket dm-event.service blk-availability.service lvm2-pvscan@.service
	fi
fi
