# OK8MP-linux-fs/rootfs/var/lib/dpkg/info/mdadm.list
/.
/etc
/etc/mdadm.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/mdadm-grow-continue@.service
/lib/systemd/system/mdadm-last-resort@.service
/lib/systemd/system/mdadm-last-resort@.timer
/lib/systemd/system/mdmon@.service
/lib/systemd/system/mdmonitor.service
/lib/systemd/system-preset
/lib/systemd/system-preset/98-mdadm.preset
/lib/systemd/system-shutdown
/lib/systemd/system-shutdown/mdadm.shutdown
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/01-md-raid-creating.rules
/lib/udev/rules.d/63-md-raid-arrays.rules
/lib/udev/rules.d/64-md-raid-assembly.rules
/lib/udev/rules.d/69-md-clustered-confirm-device.rules
/sbin
/sbin/mdadm
/sbin/mdmon
