# OK8MP-linux-fs/rootfs/var/lib/dpkg/info/nfs-utils-client.list
/.
/etc
/etc/nfs-utils.conf
/etc/nfsmount.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/nfs-statd.service
/lib/systemd/system-preset
/lib/systemd/system-preset/98-nfs-utils-client.preset
/usr
/usr/sbin
/usr/sbin/nfsstat
/usr/sbin/rpc.statd
/usr/sbin/showmount
/usr/sbin/sm-notify
/usr/sbin/start-statd
/var
/var/lib
/var/lib/nfs
/var/lib/nfs/etab
/var/lib/nfs/rmtab
/var/lib/nfs/statd
/var/lib/nfs/statd/sm
/var/lib/nfs/statd/sm.bak
/var/lib/nfs/statd/state
/var/lib/nfs/v4recovery
