# OK8MP-linux-fs/rootfs/var/lib/dpkg/info/nfs-utils-client.prerm
#!/bin/sh
[ "$1" != "upgrade" ] || exit 0
set -e
if type systemctl >/dev/null 2>/dev/null; then
	if [ -z "$D" ]; then
		systemctl stop nfs-statd.service

		systemctl disable nfs-statd.service
	fi
fi
