# OK8MP-linux-fs/rootfs/var/lib/dpkg/info/ptpd.postinst
#!/bin/sh
set -e
if type systemctl >/dev/null 2>/dev/null; then
	OPTS=""

	if [ -n "$D" ]; then
		OPTS="--root=$D"
	fi

	if [ "disable" = "enable" ]; then
		for service in ptpd.service; do
			case "${service}" in
			*@*)
				systemctl ${OPTS} enable "${service}"
				;;
			esac
		done
	fi

	if [ -z "$D" ]; then
		systemctl daemon-reload
		systemctl preset ptpd.service

		if [ "disable" = "enable" ]; then
			systemctl --no-block restart ptpd.service
		fi
	fi
fi
