# OK8MP-linux-fs/rootfs/var/lib/dpkg/info/run-postinsts.postinst
#!/bin/sh
set -e
if type systemctl >/dev/null 2>/dev/null; then
	OPTS=""

	if [ -n "$D" ]; then
		OPTS="--root=$D"
	fi

	if [ "enable" = "enable" ]; then
		for service in run-postinsts.service; do
			case "${service}" in
			*@*)
				systemctl ${OPTS} enable "${service}"
				;;
			esac
		done
	fi

	if [ -z "$D" ]; then
		systemctl daemon-reload
		systemctl preset run-postinsts.service

		if [ "enable" = "enable" ]; then
			systemctl --no-block restart run-postinsts.service
		fi
	fi
fi
