# OK8MP-linux-fs/rootfs/var/lib/dpkg/info/systemd-extra-utils.list
/.
/bin
/bin/systemd-ask-password
/bin/systemd-escape
/bin/systemd-inhibit
/bin/systemd-tty-ask-password-agent
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/multi-user.target.wants
/lib/systemd/system/sleep.target
/lib/systemd/system/sockets.target.wants
/lib/systemd/system/sysinit.target.wants
/lib/systemd/system/systemd-ask-password-console.path
/lib/systemd/system/systemd-ask-password-console.service
/lib/systemd/system/systemd-ask-password-wall.path
/lib/systemd/system/systemd-ask-password-wall.service
/lib/systemd/system/systemd-hibernate.service
/lib/systemd/system/systemd-hybrid-sleep.service
/lib/systemd/system/systemd-initctl.service
/lib/systemd/system/systemd-initctl.socket
/lib/systemd/system/systemd-suspend.service
/lib/systemd/system-generators
/lib/systemd/system-generators/systemd-gpt-auto-generator
/lib/systemd/system-sleep
/lib/systemd/systemd-ac-power
/lib/systemd/systemd-cgroups-agent
/lib/systemd/systemd-initctl
/lib/systemd/systemd-reply-password
/lib/systemd/systemd-sleep
/lib/systemd/systemd-socket-proxyd
/usr
/usr/bin
/usr/bin/systemd-cat
/usr/bin/systemd-cgls
/usr/bin/systemd-cgtop
/usr/bin/systemd-delta
/usr/bin/systemd-detect-virt
/usr/bin/systemd-path
/usr/bin/systemd-run
/usr/bin/systemd-stdio-bridge
/lib/systemd/system/multi-user.target.wants/systemd-ask-password-wall.path
/lib/systemd/system/sockets.target.wants/systemd-initctl.socket
/lib/systemd/system/sysinit.target.wants/systemd-ask-password-console.path
