# OK8MP-linux-fs/rootfs/var/lib/dpkg/info/volatile-binds.list
/.
/etc
/etc/tmpfiles.d
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/var-volatile-cache.service
/lib/systemd/system/var-volatile-lib.service
/lib/systemd/system/var-volatile-spool.service
/lib/systemd/system/var-volatile-srv.service
/lib/systemd/system-preset
/lib/systemd/system-preset/98-volatile-binds.preset
/sbin
/sbin/mount-copybind
/var
/var/cache
/etc/tmpfiles.d/etc.conf
/etc/tmpfiles.d/home.conf
