# OK8MP-linux-fs/rootfs/var/lib/dpkg/info/vsftpd.prerm
#!/bin/sh
[ "$1" != "upgrade" ] || exit 0
set -e
if type systemctl >/dev/null 2>/dev/null; then
	if [ -z "$D" ]; then
		systemctl stop vsftpd.service

		systemctl disable vsftpd.service
	fi
fi
