# OK8MP-linux-fs/rootfs/var/lib/dpkg/info/vsftpd.list
/.
/etc
/etc/default
/etc/default/volatiles
/etc/default/volatiles/99_vsftpd
/etc/pam.d
/etc/pam.d/vsftpd
/etc/tmpfiles.d
/etc/tmpfiles.d/vsftpd.conf
/etc/vsftpd.conf
/etc/vsftpd.ftpusers
/etc/vsftpd.user_list
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/vsftpd.service
/lib/systemd/system-preset
/lib/systemd/system-preset/98-vsftpd.preset
/usr
/usr/sbin
/usr/sbin/vsftpd
