# OK8MP-linux-fs/rootfs/var/lib/dpkg/info/weston-init.list
/.
/etc
/etc/profile.d
/etc/profile.d/weston.sh
/etc/udev
/etc/udev/rules.d
/etc/udev/rules.d/71-weston-drm.rules
/etc/xdg
/etc/xdg/weston
/etc/xdg/weston/weston.ini
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/weston@.service
/lib/systemd/system-preset
/lib/systemd/system-preset/98-weston-init.preset
/usr
/usr/bin
/usr/bin/weston-start
